# systemd/pony-daemons.target
[Unit]
Description=GNU/Pony Daemons
Documentation=man:pony-daemons(8)
