# systemd/pony-daemons.target
[Unit]
Description=Arch Daemons
Documentation=man:pony-daemons(8)
