# scripts/local/rc.local
#!/bin/bash
#
# Local multi-user start-up script
#


# If you have any files in £{ETC}/rc.local.d
# that are marked as executable they will be executed
# if they contain a shebang and sourced if not.

dir="£{ETC}/rc.local.d"
if [ -d "${dir}" ]; then
    for file in $(ls -1 "${dir}"); do
	file="${dir}/${file}"
	if [ -f "${file}" ] && [ -x "${file}" ]; then
	    "${file}" &
	fi
    done
fi

wait

